# cleanup-failed-shards.cron
#!/bin/sh
exec find /dcs-ssd/unpacked -maxdepth 1 -name "newshard*" -mtime +1 -delete
