# dcs-source-backend.service
.include /lib/systemd/system/dcs-common.service

[Unit]
Description=Debian Code Search: source backend
After=postgresql.service

[Service]
ExecStart=/usr/bin/source-backend -unpacked_path=/dcs/unpacked/

[Install]
WantedBy=multi-user.target
