# systemd/dcs-compute-ranking.timer
[Unit]
Description=Update ranking every hour

[Timer]
OnActiveSec=1min
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
